# infra/supervisor/rms-worker.conf
[program:rms-worker]
; Supervisor configuration for the Python dispatch engine
; Place in /etc/supervisor/conf.d/rms-worker.conf
; Apply with: supervisorctl reread && supervisorctl update

command=/path/to/rms/.venv/bin/python apps/worker/main.py
directory=/path/to/rms
autostart=true
autorestart=true
startretries=10
startsecs=5
stopwaitsecs=30

; Logging
stdout_logfile=/var/log/rms/worker-out.log
stderr_logfile=/var/log/rms/worker-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5

; Environment — supplement with production .env values
environment=NODE_ENV="production"
